# .local/sv/gdrive/finish
#!/bin/sh

exec fusermount -u /media/music
